# systemd.service.jinja
[Unit]
Description=Fenestra

[Service]
ExecStart={{ python }} {{ script_path }}
Environment=PATH={{ path }}
WorkingDirectory={{ script_folder }}
KillMode=process

[Install]
WantedBy=default.target
